# config/files/usr/lib/systemd/system/bluetooth.service.d/10-experimental.conf
[Service]
ExecStart=
ExecStart=/usr/libexec/bluetooth/bluetoothd -E d4992530-b9ec-469f-ab01-6c481c47da1c,671b10b5-42c0-4696-9227-eb28d1b049d6,330859bc-7506-492d-9370-9a6f0614037f,a6695ace-ee7f-4fb9-881a-5fac66c629af
