# files/files/usr/lib/systemd/user/profile-cleaner@.timer
[Unit]
Description=Daily Profile Cleaner on %f

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target