# files/files/usr/lib/systemd/user/profile-cleaner@.service
[Unit]
Description=Profile Cleaner on %f

[Service]
ExecStart=/usr/bin/profile-cleaner p %f