# files/files/usr/lib/systemd/system/expressvpn.service.d/override.conf
[Service]
ExecStop=/usr/bin/expressvpn disconnect