# files/files/usr/lib/systemd/system/plymouth-quit.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/plymouth quit --retain-splash
