# config/files/usr/lib/systemd/user/appman-update.service
[Unit]
Description=Auto-update Appman binary
After=local-fs.target
After=network-online.target
ConditionPathExists=%h/.local/bin/appman

[Service]
Type=oneshot
WorkingDirectory=%h
ExecStart=/usr/bin/zsh --interactive -c "appman -u"