# config/files/usr/lib/systemd/user/userbox.target
[Unit]
Description=userbox target for userbox quadlet

[Install]
WantedBy=default.target
