# config/files/usr/etc/systemd/system/gdm.service.d/override.conf
[Service]
ExecStartPre=
ExecStartPre=/bin/bash -c 'if [ -f /home/nhubao/.config/monitors.xml ]; then /usr/bin/cp /home/nhubao/.config/monitors.xml /var/lib/gdm/.config/monitors.xml; fi'
