# config/files/usr/lib/systemd/system/opt.mount
[Unit]
Description=Mount /var/opt-root to /opt

[Mount]
What=/var/opt-root
Where=/opt
Type=none
Options=bind

[Install]
WantedBy=multi-user.target
