# files/files/usr/lib/systemd/system/mdatp.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/mount -t overlay overlay -o lowerdir=/usr/lib/opt/microsoft,upperdir=/var/opt/microsoft,workdir=/var/microsoft-workdir /var/opt/microsoft
ExecStop=/usr/bin/umount /var/opt/microsoft