# files/files/usr/lib/systemd/system/mdatp.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/bash -c "env LD_LIBRARY_PATH='' /usr/bin/mount -t overlay overlay -o lowerdir=/usr/lib/opt/microsoft,upperdir=/var/microsoft-upper,workdir=/var/microsoft-workdir /var/opt/microsoft"
ExecStartPre=/usr/bin/bash -c "env LD_LIBRARY_PATH='' /usr/sbin/semodule -i /var/opt/microsoft/mdatp/conf/selinux_policies/out/audisp_mdatp.pp || true"
ExecStartPre=/usr/bin/bash -c "env LD_LIBRARY_PATH='' /usr/sbin/semanage fcontext -a -e /opt/microsoft/mdatp /var/opt/microsoft/mdatp || true"
ExecStartPre=/usr/bin/bash -c "env LD_LIBRARY_PATH='' /usr/sbin/restorecon -vR /var/opt/microsoft/mdatp"
ExecStop=/usr/bin/umount /var/opt/microsoft