# config/files/usr/etc/profile.d/askpass.sh
#!/usr/bin/env bash

export SUDO_ASKPASS='/usr/libexec/openssh/gnome-ssh-askpass'
