# files/files/usr/lib/systemd/system/flatpak-system-update.timer.d/override.conf
[Timer]
OnBootSec=2min
OnUnitActiveSec=3h
Persistent=true