# files/files/usr/lib/systemd/system/ublue-update.timer.d/override.conf
[Timer]
OnBootSec=2min
OnUnitInactiveSec=3h