# files/files/usr/lib/systemd/user/org.gnome.Shell@wayland.service.d/override.conf
[Service]
MemorySwapMax=0
MemoryZSwapMax=0